# octavaz.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/octavaz.sock

[Install]
WantedBy=sockets.target